# target/services/rpyc.service
[Unit]
Description=Target rpyc
After=network.target

[Service]
ExecStart=/home/eslab/symphony/target/run.sh
Restart=always
User=eslab
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rpycService

[Install]
WantedBy=multi-user.target

