# target/uniserver-healthlog.service
[Unit]
Description=Healthlog daemon
After=uniserver.service

[Service]
Type=forking
ExecStart=/usr/bin/bash -c "ntpdate -u pool.ntp.org && /usr/sbin/healthlogd"
ExecStop=/usr/sbin/healthlogd stop

[Install]
WantedBy=default.target
